# templates/ansible-pull-initial.service.j2
[Unit]
Description={{ job.name }} - initial setup
Before=display-manager.service getty.target ansible-pull.service
After=network-online.target
ConditionPathExists=!{{ ansible_pull_initial_file }}

[Service]
Type=oneshot
ExecStartPre=-+/bin/plymouth display-message --text="{{ ansible_pull_initial_plymouth_messsage }}"
ExecStart={{ job.job }}
ExecStopPost=-+/bin/plymouth hide-message --text="{{ ansible_pull_initial_plymouth_messsage }}"
User={{ ansible_pull_user }}

[Install]
WantedBy=first-boot-complete.target
