# templates/ansible-pull.service.j2
[Unit]
Description={{ job.name }}

[Service]
ExecStart={{ job.job }}
Type=oneshot
