# templates/ansible-pull-initial.service.j2
[Unit]
Description={{ job.name }} - initial setup
Before=display-manager.service getty@*.service ansible-pull.service
ConditionFirstBoot=yes

[Service]
Type=oneshot
ExecStartPre=-/bin/plymouth display-message --text="{{ ansible_pull_initial_plymouth_messsage }}"
ExecStart={{ job.job }}
ExecStartPost=-/bin/plymouth hide-message --text="{{ ansible_pull_initial_plymouth_messsage }}"

[Install]
WantedBy=first-boot-complete.target