# templates/ansible-pull.service.j2
[Unit]
Description={{ job.name }}

[Service]
ExecStart=/bin/bash -c '{{ job.job | regex_replace("\'", "\\'") }}'
Type=oneshot
